# util/mount/res/usb-mount@.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/bin/usb_mount.sh %i
ExecStop=/usr/local/bin/usb_umount.sh %i
